# templates/default/neutron-server.conf.erb
# Modified by Chef.

description "neutron neutron-server"
author "Hendrik Volkmer <h.volkmer@cloudbau.de>"

start on runlevel [2345]
stop on runlevel [016]

chdir /var/lib

pre-start script
  mkdir -p /var/log/neutron
  mkdir -p /var/lib/neutron
  chown neutron:root /var/log/neutron
  chown neutron:root /var/lib/neutron
end script

exec start-stop-daemon --start --chuid neutron --exec /opt/cloudbau/neutron-virtualenv/bin/neutron-server -- --config-file=/etc/neutron/neutron.conf  --log-file=/var/log/neutron/neutron-server.log --config-file=<%= @plugin_conffile_path.join(" --config-file=") %>
